# crontab
@hourly wget -O - https://kofi.karna.ge/webhook.php >/dev/null 2>&1